# deploy/ai.e3d.qwen-adapter-window-runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.e3d.qwen-adapter-window-runner</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/mini/e3d-maps/deploy/qwen_adapter_window_runner.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>/Users/mini/e3d-maps</string>
  <key>StandardOutPath</key>
  <string>/Users/mini/e3d-maps/deploy/logs/qwen-adapter-window-runner.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/mini/e3d-maps/deploy/logs/qwen-adapter-window-runner.err</string>
  <key>ThrottleInterval</key>
  <integer>30</integer>
</dict>
</plist>
