# deploy/ai.e3d.maps-runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.e3d.maps-runner</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>set -a; if [ -f /Users/mini/e3d-maps/deploy/ai.e3d.maps-runner.env ]; then source /Users/mini/e3d-maps/deploy/ai.e3d.maps-runner.env; elif [ -f /Users/mini/e3d-maps/.env ]; then source /Users/mini/e3d-maps/.env; fi; exec /opt/homebrew/bin/python3 -m agents.scheduler --loop</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mini/e3d-maps</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <!-- Restart automatically if the process exits -->
    <key>KeepAlive</key>
    <true/>

    <!-- Start at login -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/mini/e3d-maps/deploy/logs/maps-runner.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mini/e3d-maps/deploy/logs/maps-runner.err</string>

    <!-- Throttle rapid crash-restart loops -->
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
